# build_files/post-deploy-setup/post-deploy-setup.service
[Unit]
Description=LinuxBook Air per-user initial setup
ConditionPathExists=!%h/.config/linuxbook-air/initial-setup-done
ConditionPathExists=!%h/.config/linuxbook-air/initial-setup-skipped
After=graphical-session.target
PartOf=graphical-session.target
CollectMode=inactive-or-failed

[Service]
Type=simple
ExecStart=/usr/libexec/linuxbook-air-post-deploy-setup-launch.sh
# The interactive installer may wait for sudo credentials or user answers.
TimeoutStartSec=0
TimeoutStopSec=5
Restart=no

[Install]
WantedBy=graphical-session.target
