# build_files/thunderbolt-extension/linuxbook-air-thunderbolt-hotplug.service
[Unit]
Description=Keep an externally connected Thunderbolt device active

[Service]
Type=oneshot
ExecStart=/usr/libexec/linuxbook-air-thunderbolt-control hotplug-activate
TimeoutStartSec=45
