# build_files/thunderbolt-extension/linuxbook-air-thunderbolt-sleep.service
[Unit]
Description=Fully power down Thunderbolt before system sleep
Before=sleep.target
PartOf=sleep.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
ExecStart=/usr/libexec/linuxbook-air-thunderbolt-control prepare-suspend
ExecStop=/usr/libexec/linuxbook-air-thunderbolt-control resume
RemainAfterExit=yes
TimeoutStartSec=60
TimeoutStopSec=45

[Install]
RequiredBy=sleep.target
