# build_files/thunderbolt-extension/silverletter-thunderbolt-hotplug.service
[Unit]
Description=Keep an externally connected Thunderbolt device active

[Service]
Type=oneshot
ExecStart=/usr/libexec/silverletter-thunderbolt-control hotplug-activate
# Includes the optional debug discovery window and a guarded fallback teardown.
TimeoutStartSec=240
