# build_files/thunderbolt-extension/silverletter-thunderbolt-teardown.service
[Unit]
Description=Gracefully power down Thunderbolt before system sleep
Before=sleep.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/silverletter-thunderbolt-control prepare-suspend
TimeoutStartSec=60

[Install]
RequiredBy=sleep.target
