# build_files/thunderbolt-extension/silverletter-thunderbolt-hotplug.service
[Unit]
Description=Keep an externally connected Thunderbolt device active

[Service]
Type=oneshot
ExecStart=/usr/libexec/silverletter-thunderbolt-control hotplug-activate
TimeoutStartSec=45
