# build_files/thunderbolt-extension/silverletter-thunderbolt-hotplug.service
[Unit]
Description=Keep an externally connected Thunderbolt device active

[Service]
Type=oneshot
ExecStart=/usr/libexec/silverletter-thunderbolt-control hotplug-activate
# The observation window only distinguishes pending from active. It never
# invokes PCI teardown; teardown is reserved for a confirmed disconnect.
TimeoutStartSec=60
