# build_files/power/99-thunderbolt-pm.rules
# /usr/lib/udev/rules.d/99-thunderbolt-pm.rules

# Thunderbolt 2 Falcon Ridge power management (MacBook Air 7,1)
# Leave the Falcon Ridge bridge fabric and external functions entirely under
# the kernel's existing policy. This feature only changes the 07:00 NHI.
ACTION=="add", SUBSYSTEM=="pci", KERNEL=="0000:07:00.0", \
    ATTR{vendor}=="0x8086", ATTR{device}=="0x156a", \
    TEST!="/run/silverletter/thunderbolt-enabled", \
    TAG+="systemd", ENV{SYSTEMD_WANTS}+="silverletter-thunderbolt-powerdown.service"

ACTION=="add", SUBSYSTEM=="pci", KERNEL=="0000:07:00.0", \
    ATTR{vendor}=="0x8086", ATTR{device}=="0x156a", \
    TEST=="/run/silverletter/thunderbolt-enabled", ATTR{power/control}="on"

# KERNELS matches a parent, so this excludes the built-in 06:03.0 bridge itself
# and matches any external PCIe-tunnelled bridge or endpoint below it. The
# direct RUN helper only creates a cancellation marker; systemd does the slow
# driver and runtime-PM work outside the udev worker. Repeated add events also
# refresh the disconnect quiet window.
ACTION=="add", SUBSYSTEM=="pci", KERNELS=="0000:06:03.0", \
    RUN+="/usr/libexec/silverletter-thunderbolt-control hotplug-claim", \
    TAG+="systemd", ENV{SYSTEMD_WANTS}+="silverletter-thunderbolt-hotplug.service"

# Each external function emits a remove event during unplug or transient link
# training. The direct helper records activity; the delayed service only tears
# down a previously established session after a quiet interval with no add or
# remove events. DEVPATH remains available even though the removed device's
# sysfs attributes may already be gone.
ACTION=="remove", SUBSYSTEM=="pci", DEVPATH=="*/0000:06:03.0/*", \
    RUN+="/usr/libexec/silverletter-thunderbolt-control hotplug-release-request"
