# build_files/power/99-thunderbolt-pm.rules
# /usr/lib/udev/rules.d/99-thunderbolt-pm.rules

# Thunderbolt 2 Falcon Ridge power management (MacBook Air 7,1)
# Matches all TB2 bridge (156b) and NHI (156a) devices. The helper is locked
# so only one instance runs despite multiple add events.
ACTION=="add", SUBSYSTEM=="pci", ATTR{vendor}=="0x8086", ATTR{device}=="0x156b", \
    RUN+="/usr/libexec/tb-powerdown.sh"

ACTION=="add", SUBSYSTEM=="pci", ATTR{vendor}=="0x8086", ATTR{device}=="0x156a", \
    RUN+="/usr/libexec/tb-powerdown.sh"
