# build_files/power/99-thunderbolt-pm.rules
# /usr/lib/udev/rules.d/99-thunderbolt-pm.rules

# Thunderbolt 2 Falcon Ridge power management (MacBook Air 7,1)
# Matches all TB2 bridge (156b) and NHI (156a) devices.
# Teardown can block while PCI bridges disappear, so ask systemd to run the
# oneshot service instead of doing long-running work inside the udev worker.

ACTION=="add", SUBSYSTEM=="pci", ATTR{vendor}=="0x8086", ATTR{device}=="0x156b", \
    TEST!="/run/linuxbook-air/thunderbolt-enabled", \
    TAG+="systemd", ENV{SYSTEMD_WANTS}+="linuxbook-air-thunderbolt-powerdown.service"

ACTION=="add", SUBSYSTEM=="pci", ATTR{vendor}=="0x8086", ATTR{device}=="0x156a", \
    TEST!="/run/linuxbook-air/thunderbolt-enabled", \
    TAG+="systemd", ENV{SYSTEMD_WANTS}+="linuxbook-air-thunderbolt-powerdown.service"

# Enable power management for Thunderbolt controllers and ports
#ACTION=="add", SUBSYSTEM=="pci", KERNEL=="0000:05:00.0", ATTR{power/control}="auto"
#ACTION=="add", SUBSYSTEM=="pci", KERNEL=="0000:06:00.0", ATTR{power/control}="auto"
#ACTION=="add", SUBSYSTEM=="pci", KERNEL=="0000:06:03.0", ATTR{power/control}="auto"
#ACTION=="add", SUBSYSTEM=="pci", KERNEL=="0000:06:04.0", ATTR{power/control}="auto"
#ACTION=="add", SUBSYSTEM=="pci", KERNEL=="0000:06:05.0", ATTR{power/control}="auto"
#ACTION=="add", SUBSYSTEM=="pci", KERNEL=="0000:07:00.0", ATTR{power/control}="auto"

#ACTION=="add", SUBSYSTEM=="pci", ATTR{vendor}=="0x8086", ATTR{device}=="156c", \
#    RUN+="/bin/sh -c ' \
#        # 1. Target the Leaf Node (NHI) first
#        if [ -e /sys/bus/pci/devices/0000:07:00.0 ]; then \
#            echo auto > /sys/bus/pci/devices/0000:07:00.0/power/control; \
#            echo 0 > /sys/bus/pci/devices/0000:07:00.0/power/autosuspend_delay_ms; \
#        fi; \
#        # 2. Target the Bridges in reverse order (06 then 05)
#        for dev in 06:06.0 06:05.0 06:04.0 06:03.0 06:00.0 05:00.0; do \
#            if [ -e /sys/bus/pci/devices/0000:$$dev ]; then \
#                echo auto > /sys/bus/pci/devices/0000:$$dev/power/control; \
#                echo 0 > /sys/bus/pci/devices/0000:$$dev/power/autosuspend_delay_ms; \
#            fi; \
#        done; \
#        # 3. Now that they are all in D3, remove them so they stay that way
#        # and dont block the next system resume.
#        for dev in 07:00.0 06:06.0 06:05.0 06:04.0 06:03.0 06:00.0 05:00.0; do \
#            if [ -e /sys/bus/pci/devices/0000:$$dev ]; then \
#                echo 1 > /sys/bus/pci/devices/0000:$$dev/remove; \
#            fi; \
#        done; \
#    '"

    
