# build_files/power/silverletter-thunderbolt-powerdown.service
[Unit]
Description=Power down the MacBook Air Thunderbolt PCI hierarchy
After=systemd-udev-trigger.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/tb-powerdown.sh
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
